# pi-backup.cron
#!/bin/bash


/usr/local/pi-backup/run-last-pi-backup.sh --update
